# scripts/pkg/postinstall
#!/bin/bash
# Script exécuté par macOS après l'installation du .pkg (tourne en root)
set -e

# ── Trouver l'utilisateur réel (pas root) ─────────────────────────────────────
REAL_USER=$(stat -f '%Su' /dev/console 2>/dev/null || echo "$SUDO_USER")
if [ -z "$REAL_USER" ] || [ "$REAL_USER" = "root" ]; then
    REAL_USER=$(ls -l /dev/console | awk '{print $3}')
fi
REAL_HOME=$(eval echo "~$REAL_USER")
USER_UID=$(id -u "$REAL_USER")

echo "[OSMOzzz] Installation pour l'utilisateur: $REAL_USER"

# ── 1. Créer le LaunchAgent (démarrage automatique au login) ──────────────────
LAUNCH_AGENTS="$REAL_HOME/Library/LaunchAgents"
sudo -u "$REAL_USER" mkdir -p "$LAUNCH_AGENTS"
PLIST="$LAUNCH_AGENTS/com.osmozzz.daemon.plist"

cat > "$PLIST" << EOF
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.osmozzz.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/osmozzz</string>
        <string>daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>$REAL_HOME/.osmozzz/daemon.log</string>
    <key>StandardErrorPath</key>
    <string>$REAL_HOME/.osmozzz/daemon.log</string>
</dict>
</plist>
EOF

chown "$REAL_USER" "$PLIST"
echo "[OSMOzzz] LaunchAgent créé"

# ── 2. Démarrer le daemon maintenant ──────────────────────────────────────────
# macOS Ventura+ : bootstrap/bootout remplacent load/unload
launchctl bootout "gui/$USER_UID/com.osmozzz.daemon" 2>/dev/null || true
launchctl bootstrap "gui/$USER_UID" "$PLIST" 2>/dev/null || \
launchctl asuser "$USER_UID" launchctl load -w "$PLIST" 2>/dev/null || true
echo "[OSMOzzz] Daemon démarré"

# ── 3. Ouvrir le dashboard dans le navigateur (3s pour laisser le daemon s'init)
sleep 3
launchctl asuser "$USER_UID" sudo -u "$REAL_USER" open "http://localhost:7878" 2>/dev/null || \
sudo -u "$REAL_USER" open "http://localhost:7878" 2>/dev/null || true

echo "[OSMOzzz] Installation terminée — dashboard: http://localhost:7878"
